# crontab
* * * * * /root/crashdetect.sh
1,31 * * * * /root/kill.sh
